# shutdown/poweroff-control.service
[Unit]
Description=Poweroff controller
[Service]
ExecStart=/usr/bin/python /root/mimamori-plus/shutdown/pyShutdown.py
Restart=always
Type=simple
[Install]
WantedBy=multi-user.target
